# src/noodle/systemd/noodle-digest.timer
[Unit]
Description=Noodle daily digest timer

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
