# src/noodle/systemd/noodle-weekly.service
[Unit]
Description=Generate Noodle weekly review

[Service]
Type=oneshot
ExecStart=%h/.local/bin/noodle weekly
Environment="PATH=%h/.local/bin:/usr/bin"
StandardOutput=journal
