# src/noodle/systemd/noodle-weekly.timer
[Unit]
Description=Noodle weekly review timer

[Timer]
OnCalendar=Sun *-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
